# data/com.vysp3r.ProtonPlus.timer
[Unit]
Description=ProtonPlus Notification Timer

[Timer]
OnBootSec=0
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
