# data/com.vysp3r.ProtonPlus.service
[Unit]
Description=ProtonPlus Notification Service

[Service]
Type=oneshot
ExecStart=/usr/bin/notify-send "Hi!"
